# bootstrap/gunicorn_supervisor.conf
[program:wormhole_gunicorn_live]
command=/home/pi/.virtualenvs/wormhole/bin/gunicorn wormhole_web:app -b 127.0.0.1:1559 -u root
directory=/home/pi/Devel/wormhole
user=root
autostart=True
autorestart=True
redirect_stderr=True
